# config/supervisor/project_name_celerybeat.conf
[program:{{ project_name }}_celerybeat]
command=/home/{{ project_name }}/{{ project_name }}/scripts/celerybeat.bash
user={{ project_name }}
