# config/supervisor/project_name.conf
[program:{{ project_name }}]
command = /home/{{ project_name }}/{{ project_name }}/scripts/runserver.bash
user = {{ project_name }}
