# config/supervisor/splitkart_celery.conf
[program:{{ project_name }}_celery]
command=/home/{{ project_name }}/{{ project_name }}/scripts/celery.bash
user={{ project_name }}
